# ebook-generator.service
[Unit]
Description=Ebook Generator (Streamlit :8501 + API :8765)
After=network.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw/projects/1ai-ebook
Environment=PATH=/home/openclaw/.local/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
Environment=OMNIROUTE_BASE_URL=http://localhost:20128/v1
Environment=EBOOK_API_KEY=dev-key-change-me
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3 run.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
